# deployment/files/etc/supervisord/conf.d/voternet.conf
[program:{{ project_name }}]
command={{project_root }}/code/run.fcgi {{app_port}}
user={{ project_name }}
directory={{ project_root }}/code
environment=PATH="{{ project_root }}/env/bin"
redirect_stderr=true
stdout_logfile={{ project_root }}/logs/{{project_name}}.log